# esp-9010_lmp_esw_v00_06/config/includes.chroot/etc/profile.d/01extra_path.sh
#!/bin/bash

export PATH="$PATH:/opt/test/mfg"
